# cli_tools/gce_inventory_agent/google-compute-engine-inventory.conf
# Google Compute Engine inventory agent
start on stopped google-compute-engine-inventory

respawn
exec /usr/bin/gce_inventory_agent