# cli_tools/google-osconfig-agent/google-osconfig-agent.conf
description "Google OSConfig Agent"
start on (started networking)

expect fork
respawn
exec /usr/bin/google_osconfig_agent
