# cli_tools/google-osconfig-agent/google-osconfig-agent.service
[Unit]
Description=Google OSConfig Agent

[Service]
ExecStart=/usr/bin/google_osconfig_agent
Restart=always
RestartSec=120
