# cli_tools/google-osconfig-agent/google-osconfig-agent.conf
description "Google OSConfig Agent"
start on stopped rc RUNLEVEL=[2345]

respawn
exec /usr/bin/google_osconfig_agent
