# daisy_workflows/build-publish/rhui/cds_artifacts/rhui-health-check.timer
[Unit]
Description=Timer for CDS health check
Requires=rhui-health-check.service

[Timer]
Unit=rhui-health-check.service
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
